# target/device/peaslaker/ubit/target_skeleton/etc/network/if-up.d/ntpd
#!/bin/sh

[ "$IFACE" = "eth0" ] || exit 0

/usr/sbin/ntpd -q -p pool.ntp.org
